# ops/systemd/glitch-signal.service
[Unit]
Description=Glitch Social Media Agent — autonomous social video + ORM agent
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=support
Group=support
WorkingDirectory=/home/support/glitch-signal
EnvironmentFile=/home/support/glitch-signal/.env
ExecStart=/home/support/glitch-signal/.venv/bin/uvicorn glitch_signal.server:app --host 127.0.0.1 --port 3111
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=glitch-signal

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/support/glitch-signal /var/lib/glitch-signal

[Install]
WantedBy=multi-user.target
